# installation/debian/postinst
#!/bin/sh
set -e

#DEBHELPER#

systemctl daemon-reload
systemctl reenable touchegg
systemctl stop touchegg
systemctl start touchegg

exit 0
